# docker/buildbot-slave/supervisord.conf
[supervisord]
nodaemon=true

[program:khpy]
command=/opt/khpy/khs restart

[program:buildslave]
command=buildslave start --nodaemon slave
environment=HOME="/home/buildbot"
